# roles/domain/templates/noip2.service.j2
[Unit]
Description=No-IP Dynamic DNS Update Service
After=network.target

[Service]
Type=forking
ExecStart=/usr/local/bin/noip2
Restart=always

[Install]
WantedBy=multi-user.target